# ops/systemd/worldcup-backup.service
[Unit]
Description=Backup World Cup Collector Data

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/backup-worldcup-collector
